# deploy/oko.service
[Unit]
Description=Oko Discord scam-detection bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oko
WorkingDirectory=/opt/oko
ExecStart=/opt/oko/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
